# deploy/vnet-backend.service
[Unit]
Description=VNet Backend - Spring Boot Application
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/vnet/backend

# Load Environment Variables from external .env file (Not committed to Git)
EnvironmentFile=/var/www/vnet/backend/.env

# Jalankan aplikasi
ExecStart=/usr/bin/java -jar /var/www/vnet/backend/vnet-backend-0.0.1-SNAPSHOT.jar

# Restart otomatis kalau crash
Restart=always
RestartSec=10

# Log output
StandardOutput=append:/var/www/vnet/logs/backend.log
StandardError=append:/var/www/vnet/logs/backend-error.log

[Install]
WantedBy=multi-user.target
